# descriptions/suse-leap-15.1/root/etc/profile.d/sesame.sh
export PATH=$PATH:/opt/sesame/bin:/opt/sesame/sbin
